# lxc/upstart/authentication.conf
start on filesystem and net-device-up IFACE!=lo
script
    java -jar /authentication/AuthenticationApplication.jar server /authentication/config.yml
end script
console none
respawn
