# lxc/upstart/session.conf
start on filesystem and net-device-up IFACE!=lo
script
    java -jar /session/SessionApplication.jar server /session/config.yml
end script
console none
respawn
